# clipmesh.service
[Unit]
Description=clipmesh - encrypted LAN clipboard sync
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/.cargo/bin/clipmesh
# always: a config-file change makes clipmesh exit cleanly to reload, and any
# unexpected exit should bring it back too.
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
